# cron/april_fool.cron
*	4	1	7	30	/home/pi/RasbellyPi/cron/randbellypi

