# cron/april_fool.cron
*	4	1	7	0	/home/pi/RasbellyPi/cron/randbellypi